# lab09-power-consumption-monitoring-with-powertop/scripts/auto-power-mgmt.service
[Unit]
Description=Automated Power Management Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/auto_power_mgmt.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
