# lab09-power-consumption-monitoring-with-powertop/scripts/power-optimize.service
[Unit]
Description=Power Optimization Service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-optimize.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
